# launchd/com.xiongbatian.12-news-radar.weekly-plan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xiongbatian.12-news-radar.weekly-plan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/xiongbatian/Documents/Codex/2026-07-04/ai-builder/work/12-news-radar/scripts/news_radar.py</string>
    <string>weekly-plan</string>
    <string>--publish</string>
    <string>--notify</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/xiongbatian/Documents/Codex/2026-07-04/ai-builder/work/12-news-radar</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>40</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/12-news-radar.weekly-plan.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/12-news-radar.weekly-plan.err.log</string>
</dict>
</plist>
